# aura_backend/aura-backend.service
[Unit]
Description=Aura Advanced AI Companion Backend
After=network.target

[Service]
Type=simple
User=ty
WorkingDirectory=/home/ty/Repositories/ai_workspace/emotion_ai/aura_backend
Environment=PATH=/home/ty/Repositories/ai_workspace/emotion_ai/aura_backend/venv/bin:/usr/bin:/bin
ExecStart=/home/ty/Repositories/ai_workspace/emotion_ai/aura_backend/venv/bin/python /home/ty/Repositories/ai_workspace/emotion_ai/aura_backend/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
